# dnsforwarder.service.template
[Unit]
Description=Fix for Telcel DNS blocking
BindsTo=usbtether.service
After=usbtether.service
[Service]
ExecStart=$PWD/forward.py
Restart=always
RestartSec=10
[Install]
WantedBy=multi-user.target
